# sloelux-perfbot/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fastapi]
command=uvicorn fastapi_server:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/fastapi.err.log
stdout_logfile=/app/logs/fastapi.out.log
environment=PYTHONPATH="/app"

[program:mcp_server]
command=python mcp_server.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/mcp.err.log
stdout_logfile=/app/logs/mcp.out.log
environment=PYTHONPATH="/app"

[program:perf_loop]
command=python perf_loop.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/perf_loop.err.log
stdout_logfile=/app/logs/perf_loop.out.log
environment=PYTHONPATH="/app" 